# scripts/systemd/gold-standard-daily.timer
[Unit]
Description=Daily timer to run Syndicate main one-shot at 05:00

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
