# deploy/systemd/gold-standard-retry.service
[Unit]
Description=Gold Standard Notion Retry Service
After=network.target

[Service]
Type=oneshot
User=ali
ExecStart=/usr/bin/env python /mnt/disk/gold_standard/scripts/retry_failed_publishes.py

[Install]
WantedBy=multi-user.target
