# deploy/systemd/gold-standard-model-cleanup.timer
[Unit]
Description=Run Gold Standard model cleanup weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
