# deploy/systemd/gold-standard-run-once.timer
[Unit]
Description=Description=Run Syndicate Once Daily at 5 AM

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
