# deploy/systemd/gold-standard-offloaded-executor.service
[Unit]
Description=Syndicate Offloaded Executor (local fast GGUF models)
After=network.target

[Service]
Type=simple
ExecStart=/home/adam/worxpace/.venv/bin/python /home/adam/worxpace/syndicate/scripts/offloaded_executor.py
WorkingDirectory=/home/adam/worxpace/syndicate
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
