# deploy/systemd/gold-standard-daily-llm-report.timer
[Unit]
Description=Description=Run Syndicate daily LLM report daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
