# deploy/systemd/gold-standard-premarket-watcher.service
[Unit]
Description=Syndicate PreMarket Watcher
After=network.target

[Service]
Type=simple
User=adam
WorkingDirectory=/home/adam/worxpace/syndicate
ExecStart=/home/adam/worxpace/.venv/bin/python scripts/premarket_watcher.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
