# deploy/systemd/gold-standard-model-cleanup.service
[Unit]
Description=Syndicate Model Cleanup (dry-run by default)
After=network.target

[Service]
Type=oneshot
ExecStart=/home/adam/worxpace/.venv/bin/python /home/adam/worxpace/syndicate/scripts/cleanup_models.py --confirm
WorkingDirectory=/home/adam/worxpace/syndicate
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
