# deploy/systemd/gold-standard-daily-llm-report.service
[Unit]
Description=Syndicate daily LLM operations report
After=network.target

[Service]
Type=oneshot
User=adam
WorkingDirectory=/home/adam/worxpace/syndicate
EnvironmentFile=/home/adam/worxpace/syndicate/.env
Environment=PYTHONPATH=/home/adam/worxpace/syndicate/src
ExecStart=/home/adam/worxpace/.venv/bin/python -m digest_bot.daily_report

[Install]
WantedBy=multi-user.target
