# deploy/systemd/gold-standard-llm-worker.service
[Unit]
Description=Syndicate LLM Worker
After=network.target

[Service]
Type=simple
User=adam
Group=adam
WorkingDirectory=/home/adam/worxpace/syndicate
ExecStart=/opt/run_syndicate_llm_worker.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
