# deploy/systemd/normalized/syndicate-publish-digest.service
[Unit]
Description=Syndicate - Digest Publishing Service
After=network.target

[Service]
Type=oneshot
User=ali_shakil_backup_gmail_com
WorkingDirectory=/home/ali_shakil_backup_gmail_com/syndicate
ExecStart=/home/ali_shakil_backup_gmail_com/miniforge3/envs/syndicate/bin/python scripts/scheduled_publisher.py digest
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
